# 3. Installation/3. Docker Containers/systemd/backup@.timer
[Unit]
Description=Hourly backup of %i content.

# Dependency binding
BindsTo=universal@%i.service

[Timer]
OnActiveSec=20
OnUnitActiveSec=60min

