# 3. Installation/3. Docker Containers/systemd/backup@.service
[Unit]
Description=Back up data from %i

[Service]
Type=oneshot
TimeoutStartSec=0
WorkingDirectory=/data/domains/%i/
ExecStart=/bin/bash -euxc '/data/domains/%i/BACKUP'
